# tools/supervisord.conf
[program:gunicorn]
directory=/home/yashh/bubly
command=gunicorn_django -c /home/yashh/bubly/gunicorn.conf.py
process_name=%(program_name)s
numprocs=1
priority=999
user=yashh

